# node/app885.service
[Unit]
Description=Acesso ao Banco de Consultas do Protheus
After=network.target
StartLimitIntervalSec=0
[Service]
Type=simple
Restart=always
RestartSec=1
User=root
ExecStart=/usr/bin/env nodejs /var/www/node/app885.js
[Install]
WantedBy=multi-user.target
